# docker/postgresql/supervisord-postgresql.conf
[program:postgresql]
command=/docker/start-postgresql.sh
numprocs=1
autostart=true
autorestart=true
