# docker/postgresql/supervisord-postgresql-build.conf
[program:postgresql-build]
command=/docker/start-postgresql-build.sh
numprocs=1
autostart=true
autorestart=false
