# docker/web/supervisord-system.conf
[program:system]
command=/docker/start-system.sh
numprocs=1
autostart=true
autorestart=false
