# docker/web/supervisord-rserve.conf
[program:rserve]
command=/docker/start-rserve.sh
numprocs=1
autostart=true
autorestart=true
