# docker/web/supervisord-apache2.conf
[program:apache2]
command=/docker/start-apache2.sh
numprocs=1
autostart=true
autorestart=true
