# docker/enketo/supervisord-redis-cache.conf
[program:redis-cache]
command=/docker/start-redis-cache.sh
numprocs=1
autostart=true
autorestart=true
