# docker/enketo/supervisord-redis-main.conf
[program:redis-main]
command=/docker/start-redis-main.sh
numprocs=1
autostart=true
autorestart=true
