# docker/enketo/supervisord-pm2.conf
[program:pm2]
command=/docker/start-pm2.sh
numprocs=1
autostart=true
autorestart=true
