# docker/mongo/supervisord-sshd.conf
[program:sshd]
command=/docker_add/start-sshd.sh
numprocs=1
autostart=true
autorestart=true
