# docker/mongo/supervisord-mongo.conf
[program:mongo]
command=/docker_add/start-mongo.sh
numprocs=1
autostart=true
autorestart=true
